# systemd/brokersaurusrex.socket
[Unit]
Description=Socket for Brokersaurusrex

[Socket]
ListenStream=127.0.0.1:302
Accept=yes

[Install]
WantedBy=sockets.target

